# service/podman/run
#!/bin/sh
exec chpst -e "$TURNSTILE_ENV_DIR" podman system service -t 0